# arch/install/.config/autostart/pipewire-pulse.desktop
[Desktop Entry]
Type=Application
Name=PipeWire Pulse
Exec=pipewire-pulse
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
