# arch/install/.config/autostart/pulsemixer.desktop
[Desktop Entry]
Type=Application
Name=Pulsemixer
Exec=pulsemixer
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
